# dashboard_web/careplus-dashboard.service
[Unit]
Description=CarePlus Sprint 03 FIWARE Dashboard
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
WorkingDirectory=/opt/careplus-dashboard
Environment=FIWARE_HOST=34.69.120.192
Environment=STH_COMET_URL=http://127.0.0.1:8666
Environment=DASHBOARD_PORT=8080
Environment=FIWARE_SERVICE=openiot
Environment=FIWARE_SERVICE_PATH=/
Environment=ENTITY_ID=CarePlusMission:totem001
Environment=ENTITY_TYPE=CarePlusWalkingMission
ExecStart=/opt/careplus-dashboard/venv/bin/python /opt/careplus-dashboard/app.py
Restart=always
RestartSec=5

[Install]
WantedBy=multi-user.target
